# crontabs/root
MAILTO=ecoreuse9@gmail.com

# Server management
@daily sudo reboot
0 1 * * * sudo apt-get update && sudo apt-get upgrade -y
0 2 * * * sudo apt-get autoremove -y
