# crontabs/root
MAILTO=ecoreuse9@gmail.com

# Server management
0 1 * * * sudo reboot
30 0 * * * sudo apt-upgrade.sh
