# crontabs/root
MAILTO=ecoreuse9@gmail.com

# Server management
@daily sudo reboot
0 1 * * * sudo apt-upgrade.sh
